# deploy/linux/openclaw-watchdog.service.example
[Unit]
Description=OpenClaw Watchdog
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/openclaw-watchdog
ExecStart=%t/openclaw-watchdog/node_modules/.bin/node src/cli/index.js run
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
